# com.skydiving.video-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skydiving.video-organizer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/connor/Skydiving/video_organiser/watch_videos.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/connor/Skydiving/video_organiser/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/connor/Skydiving/video_organiser/output.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/connor</string>
    </dict>
</dict>
</plist> 